# chef-cookbooks/jboss/templates/jboss_init.erb
#!/bin/sh
# chkconfig: 2345 95 05

start(){
echo "Starting jboss.."
su -l <%= @jboss_user %> -c "<%= @jboss_home %>/bin/run.sh -b <%= @jboss_ip %> & "
}
stop(){
echo "Stopping jboss.."
su -l <%= @jboss_user %> -c "<%= @jboss_home %>/bin/shutdown.sh -S & "
}
restart(){
stop
sleep 20
start
}
case "$1" in
start)
start
;;
stop)
stop
;;
restart)
restart
;;
*)
echo "Usage: jboss {start|stop|restart}"
exit 1
esac
exit 0
